# calist.service
[Unit]
Description=Calist
After=network.target

[Service]
Type=simple
User=calist
Group=calist
WorkingDirectory=/opt/calist
ExecStart=/usr/bin/env bun run src/server/index.ts
Environment=NODE_ENV=production
Environment=PORT=3100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target