# scheduler.service
[Unit]
Description=Personal Scheduler
After=network.target

[Service]
Type=simple
User=scheduler
Group=scheduler
WorkingDirectory=/opt/scheduler
ExecStart=/usr/bin/env bun run src/server/index.ts
Environment=NODE_ENV=production
Environment=PORT=3100
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target